# Millllllll/gofish-1.2/gofish/debian/maint/postrm
#! /bin/sh

set -e

deluser gopher || true

if [ "$1" = "purge" ] ; then
	update-rc.d gofish remove >/dev/null
fi
